# eye.service
[Unit]
Description = Eye agent
After network.target = auditd.service

[Service]
Type = simple
ExecStart = /opt/eye/eye-agent
