# DEBIAN/postinst
#!/bin/sh

systemctl enable eye.service
systemctl start eye.service
