# DEBIAN/prerm
#!/bin/sh

systemctl stop eye.service
systemctl disable eye.service
